# spark/supervisord.conf
[supervisord]
nodaemon=true

[program:spark]
command=/bin/bash -c "bootstrap.sh"