# hbase-onbuild/supervisord.conf
[supervisord]
nodaemon=true

[program:hbase]
command=/bin/bash -c "bootstrap.sh -DFOREGROUND"